# telegram-bot-ram.service
[Unit]
Description=Telegram Auto-Join Bot (Ram)
After=network.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/home/bot/telegram-bots/ram-tg-bot
Environment=PATH=/home/bot/telegram-bots/ram-tg-bot/venv/bin
ExecStart=/home/bot/telegram-bots/ram-tg-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
